# recipes-core/initscripts/initscripts/firstboot-build_version.service

Description=intel-aero

[Service]
Type=oneshot
ExecStart=/etc/init.d/build_version.sh

[Install]
WantedBy=multi-user.target

